# debian/postinst
#!/bin/bash -e
#
# ==============================================================================
# PAQUETE: kit-servicios
# ARCHIVO: postinst
# DESCRIPCIÓN: Configura el sistema despues la instalación del paquete.
# COPYRIGHT:
#  (C) 2016 Victor Pino <victopin0@gmail.com>
# LICENCIA: GPL3
# ==============================================================================
#
# Este programa es software libre. Puede redistribuirlo y/o modificarlo bajo los
# términos de la Licencia Pública General de GNU (versión 3).

case "$1" in
    configure)
        
        #echo "Colocando /etc/default/celeryd en /etc/kds/diverted/celeryd"
        #dpkg-divert --divert /etc/kds/diverted/celeryd --rename /etc/default/celeryd

        #echo "Instalando configuracion personalizada para Celeryd"
        #ln /etc/kds/celeryd /etc/default/celeryd

        if [ -f  /var/log/celery ] ; then

            echo "Eliminando los logs"
            rm -rf /var/log/celery
            rm -rf /var/run/celery

        fi
        
        echo "Echo Agregando usuario"
        adduser --system --shell /bin/bash --gecos 'KDS system user' --group --disabled-password --quiet --no-create-home kds

        /etc/init.d/rabbitmq-server restart

        echo "Configurando rabbitmq-server"

        rabbitmqctl add_user kds 11

        echo "==== add virtual host 'kds_vhost' === "

        rabbitmqctl add_vhost kds_vhost

        echo "==== add user tag 'kds_tag' for user 'kds' ==="

        rabbitmqctl set_user_tags kds kds_tag

        echo "==== set permission for user 'kds' on virtual host 'kds_vhost' ==="

        rabbitmqctl set_permissions -p kds_vhost kds ".*" ".*" ".*"

        echo "Eliminar archivo /etc/nginx/sites-available/default"
        if [ -f /etc/nginx/sites-available/default ] ; then

            rm -rf /etc/nginx/sites-available/default

        fi

        if [ -d /usr/share/python/kit-servicios/ ]; then
        
            echo "Modificando los permisos de la carpeta del proyecto"
            chmod 755 /usr/share/python/kit-servicios/

        fi
        
        echo "Permisos a la carpeta del proyecto"
        chown -R www-data:www-data /usr/share/python/kit-servicios/lib/python2.7/site-packages/kit-servicios/

        echo "Reiniciando los servicios"
		service nginx restart
		service uwsgi restart || true
        #/etc/init.d/celeryd restart 
        service redis-server restart
        
        #Demonio encargado de supervisar celery
        supervisorctl reread
        supervisorctl update

        #Demonio encargado de el modo emperador de uwsgi
        systemctl daemon-reload
        systemctl enable nginx emperor.uwsgi
        systemctl start nginx emperor.uwsgi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
